# app/conf.d/gost.conf
[program:gost]
command=/usr/local/bin/gost -L socks5://:10801 -L tcp://:2222/%(ENV_REMOTE_HOST)s:2222 -L http://:18888
autorestart=true
