# app/conf.d/simpleproxy.conf
[program:simpleproxy]
command=simpleproxy -L 0.0.0.0:2222 -R "%(ENV_REMOTE_HOST)s"
autorestart=true