# app/conf.d/univpn.conf
[program:univpn]
directory=/usr/local/UniVPN
command=/usr/local/UniVPN/UniVPN
autorestart=true
